# notify-reg.service.in
[Unit]
Description=Notification Registry Service
After=dbus-session.service

[Service]
Type=simple
ExecStart={{BIN_PATH}} --path %h/notify-reg.txt
Restart=always

[Install]
WantedBy=default.target
