# notify-reg.service.in
[Unit]
Description=Qutebrowser Floccus Sync Daemon
After=network-online.target

[Service]
Type=simple
ExecStart={{BIN_PATH}} --path %h/notify-reg.txt
Restart=on-failure
RestartSec=60s

[Install]
WantedBy=default.target
